# chromeos-base/chromeos-bsp-rockpi4b/files/udev-rules/99-mali-wait.rules
KERNEL=="mali0", SUBSYSTEM=="misc", ACTION=="add", RUN+="/sbin/start ui"
